# .config/systemd/user/kdeconnect.service
[Unit]
Description=KDE Connect daemon

[Service]
ExecStart=/usr/sbin/kdeconnectd
Restart=on-failure

[Install]
WantedBy=default.target
