# .config/systemd/user/nm-applet.service
[Unit]
Description=NetworkManager Applet

[Service]
ExecStart=/usr/sbin/nm-applet
Restart=on-failure

[Install]
WantedBy=default.target
