# .config/systemd/user/picom.service
[Unit]
Description=Picom compositor

[Service]
ExecStart=/usr/sbin/picom
Restart=on-failure

[Install]
WantedBy=default.target
