# .config/systemd/user/polybar.service
[Unit]
Description=Polybar status bar

[Service]
ExecStart=/home/%u/.config/polybar/launch.sh
Restart=on-failure

[Install]
WantedBy=default.target
